# templates/systemd.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Renew LE certificates

[Timer]
OnCalendar={{ certbot_periodicity }}
Unit=certbot.service

[Install]
WantedBy=multi-user.target
